# install/etc/services.d/clashx/run
#!/command/with-contenv bash

PROCESS_NAME="clashx"
cd /root/subconverter
./subconverter -g 
echo "Starting clashx"
clash -d /root/.config/clash -f /root/.config/clash/pgfast.yml
